# deploy/systemd/adc.service
[Unit]
Description=WSO2 API Discovery Controller
Documentation=https://github.com/Kidhurshan/WSO2-API-Discovery-Controller
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=adc
Group=adc
# Database credentials are loaded from /etc/adc/adc.env so they never appear
# in the config file. ADC's loader expands ${POSTGRES_USER}, ${POSTGRES_PASSWORD},
# and ${POSTGRES_DB} references in /etc/adc/config.toml against this environment.
# The leading "-" makes the file optional — installs that hard-code credentials
# in config.toml will still work.
EnvironmentFile=-/etc/adc/adc.env
ExecStart=/usr/local/bin/adc --config /etc/adc/config.toml
Restart=always
RestartSec=10
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/adc
PrivateTmp=true

[Install]
WantedBy=multi-user.target
