# deploy/systemd/adc.service
[Unit]
Description=WSO2 API Discovery Controller
Documentation=https://github.com/Kidhurshan/WSO2-API-Discovery-Controller
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=adc
Group=adc
# Deployment-mode vars with bundled-VM defaults. To use external PostgreSQL,
# set POSTGRES_HOST=your.db.host in /etc/adc/adc.env — EnvironmentFile is
# processed after Environment=, so file values override these defaults.
Environment="ADC_MODE=standalone" "POSTGRES_HOST=localhost"
# Database credentials are loaded from /etc/adc/adc.env so they never appear
# in the config file. ADC's loader expands ${POSTGRES_USER}, ${POSTGRES_PASSWORD},
# ${POSTGRES_DB}, and ${POSTGRES_HOST} references in /etc/adc/config.toml
# against this environment. The leading "-" makes the file optional — installs
# that hard-code credentials in config.toml will still work.
EnvironmentFile=-/etc/adc/adc.env
ExecStart=/usr/local/bin/adc --config /etc/adc/config.toml
Restart=always
RestartSec=10
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/adc
PrivateTmp=true

[Install]
WantedBy=multi-user.target
